# systemd/awairelement-web.service
# Dev reference copy — the canonical unit is templated by the native-apps
# Ansible role in the homelab repo (awairelement-web.service.j2).
# Restarted alongside the poller by restart.sh, so `itguy deploy awairelement`
# updates both processes.
[Unit]
Description=Awair Element dashboard (gunicorn)
After=network-online.target

[Service]
WorkingDirectory=%h/sources/awairelement
ExecStart=%h/sources/awairelement/.venv/bin/gunicorn -w 2 -b 0.0.0.0:8097 'awair.web:create_app()'
EnvironmentFile=-%h/.config/awairelement/environment
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
