# systemd/awair-poller.service
# Dev reference copy — the canonical unit is templated by the native-apps
# Ansible role in the homelab repo (slice 4). Install manually meanwhile:
#   ln -sf ~/sources/awairelement/systemd/awair-poller.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now awair-poller
[Unit]
Description=Awair Element air quality poller
After=network-online.target

[Service]
WorkingDirectory=%h/sources/awairelement
ExecStart=%h/sources/awairelement/.venv/bin/python -m awair.poller
EnvironmentFile=-%h/.config/awairelement/env
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
